# sheikha-main-portal/packaging/deb/DEBIAN/prerm
#!/bin/sh
set -e

SERVICE_NAME="sheikha"

if [ -f "/etc/systemd/system/${SERVICE_NAME}.service" ]; then
    systemctl disable --now "${SERVICE_NAME}" || true
fi

exit 0
